# launchd/com.housing-monitor.hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.housing-monitor.hourly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/macmini3/housing_monitor/scripts/run_hourly.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/macmini3/housing_monitor</string>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/macmini3/housing_monitor/logs/hourly.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/macmini3/housing_monitor/logs/hourly.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/macmini3</string>
        <key>PATH</key>
        <string>/opt/homebrew/Caskroom/miniforge/base/envs/housing/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
